# crontab
# Run the script every day at midnight
0 0 * * * /app/script.sh >> /var/log/script.log 2>&1
